# cyberkitty19-transkribator.service
[Unit]
Description=Cyberkitty19 Transkribator Bot Service
Documentation=https://github.com/your-username/cyberkitty19-transkribator
After=network.target docker.service
Requires=docker.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
User=cyberkitty
Group=cyberkitty
WorkingDirectory=/opt/cyberkitty19-transkribator
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStartPre=/usr/local/bin/docker-compose down
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
ExecReload=/usr/local/bin/docker-compose restart
TimeoutStartSec=300
TimeoutStopSec=60
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target 